# Luci-app/luci-app-adguardhome/root/etc/uci-defaults/40_luci-AdGuardHome
#!/bin/sh

normalize_binpath_value() {
	local path="$1" parent
	[ -n "$path" ] || path='/etc/config/adGuardConfig/AdGuardHome'
	while [ "${path%/}" != "$path" ]; do
		path="${path%/}"
	done
	[ -n "$path" ] || path='/etc/config/adGuardConfig/AdGuardHome'
	parent="${path%/*}"
	if [ "${path##*/}" = 'AdGuardHome' ] && [ "$parent" != "$path" ] && [ "${parent##*/}" = 'AdGuardHome' ] && [ -e "$parent" ] && [ ! -d "$parent" ]; then
		path="$parent"
	fi
	if [ -d "$path" ]; then
		printf '%s/AdGuardHome\n' "$path"
	else
		printf '%s\n' "$path"
	fi
}

configpath=$(uci -q get AdGuardHome.AdGuardHome.configpath 2>/dev/null)
workdir=$(uci -q get AdGuardHome.AdGuardHome.workdir 2>/dev/null)
binpath=$(uci -q get AdGuardHome.AdGuardHome.binpath 2>/dev/null)
downloadarch=$(uci -q get AdGuardHome.AdGuardHome.downloadarch 2>/dev/null)
release_channel=$(uci -q get AdGuardHome.AdGuardHome.release_channel 2>/dev/null)
tagname=$(uci -q get AdGuardHome.AdGuardHome.tagname 2>/dev/null)
binpath=$(normalize_binpath_value "$binpath")

uci -q batch <<-EOF >/dev/null 2>&1
	delete ucitrack.@AdGuardHome[-1]
	add ucitrack AdGuardHome
	set ucitrack.@AdGuardHome[-1].init=AdGuardHome
	commit ucitrack
	delete AdGuardHome.AdGuardHome.ucitracktest
	set AdGuardHome.AdGuardHome.configpath='${configpath:-/etc/config/adGuardConfig/AdGuardHome.yaml}'
	set AdGuardHome.AdGuardHome.workdir='${workdir:-/etc/config/adGuardConfig/workspace}'
	set AdGuardHome.AdGuardHome.binpath='${binpath:-/etc/config/adGuardConfig/AdGuardHome}'
	set AdGuardHome.AdGuardHome.downloadarch='${downloadarch:-auto}'
	set AdGuardHome.AdGuardHome.release_channel='${release_channel:-release}'
	set AdGuardHome.AdGuardHome.tagname='${tagname:-release}'
	commit AdGuardHome
EOF

mkdir -p /etc/config/adGuardConfig/workspace 2>/dev/null
rm -rf /tmp/luci-modulecache/ /var/luci-modulecache/
rm -f /tmp/luci-indexcache /tmp/luci-indexcache.* /var/luci-indexcache.* 2>/dev/null
if [ -x /etc/init.d/rpcd ]; then
	/etc/init.d/rpcd reload >/dev/null 2>&1 || true
	( sleep 1; /etc/init.d/rpcd restart >/dev/null 2>&1 ) &
fi
exit 0
